# work/70-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="su clintedwards -c \"/usr/bin/autorandr --change\""
